# ops/se.dyallo.vault-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>se.dyallo.vault-dashboard</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/Dennis.Dyall/Code/other/dyallose/scripts/regenerate-vault-dashboard.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>1800</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/Dennis.Dyall</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/Dennis.Dyall/Code/other/dyallose/Logs/vault-dashboard.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/Dennis.Dyall/Code/other/dyallose/Logs/vault-dashboard.error.log</string>
</dict>
</plist>
